# sensor-net-data-api.socket
[Socket]
ListenStream=8001

[Install]
WantedBy=sockets.target
